# build/sickchill.conf
[program:sickchill]
user = nobody
command = /usr/sbin/python2 /opt/sickchill/SickBeard.py --config /config/config.ini --datadir /config
umask = 000
