# build/sickchill.conf
[program:sickchill]
user = nobody
command = /usr/sbin/python3 /opt/sickchill/SickChill.py --config /config/config.ini --datadir /config
umask = 000

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh 'sickchill'
umask = 000
